# raspberry/systemd/aquirejson.service
[Unit]
Description=Fetch sensor data and parse it to a json file

[Service]
Type=simple
ExecStart=/usr/bin/readsensors.py

