# raspberry/configs/clock-sync.service
[Unit]
Description=Update system clock from RTC

[Service]
Type=oneshot
ExecStart=/usr/bin/hwclock -s

[Install]
WantedBy=multi-user.target
